# sandbox/rust/wcircle/dist/wcircle.service
[Unit]
Description=wcircle circular scrolling daemon
After=graphical-session.target

[Service]
ExecStart=%h/.cargo/bin/wcircle
Restart=on-failure

[Install]
WantedBy=default.target
